# root/etc/s6-overlay/s6-rc.d/svc-novnc/run
#!/usr/bin/with-contenv bash

cd $HOME
exec s6-setuidgid abc \
  /bin/bash /usr/sbin/novnc \
    --vnc "$VNC_HOST:$VNC_PORT" \
    --listen "$NOVNC_HOST:$NOVNC_PORT"
